# data/ayatana-indicator-datetime.conf.in
description "Ayatana Indicator Date & Time Backend"

start on ayatana-indicator-services-start
stop on desktop-end or ayatana-indicator-services-end

respawn
respawn limit 2 10

exec @pkglibexecdir@/ayatana-indicator-datetime-service
