# data/ayatana-indicator-datetime.service.in
[Unit]
Description=Ayatana Indicator DateTime Service
PartOf=graphical-session.target
PartOf=ayatana-indicators.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-datetime-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
